# roles/lab_debian/files/etc/profile.d/99local.sh
# overwrite PATH
export PATH="/usr/sbin:/usr/bin:/sbin:/bin"

# set custom umask
umask 777
